# pkg/aether/deploy/dice-export.timer
# systemd USER timer — install to ~/.config/systemd/user/dice-export.timer
#
# Runs dice-export.service nightly to keep /dice ~1 day fresh. Persistent so a
# missed run (host asleep) fires on next boot. See deploy/DICE.md.

[Unit]
Description=Faerrin aether — nightly dice-data export + rebuild

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
