# pkg/wretch/deploy/wretch-reconcile.service
# systemd USER unit — install to ~/.config/systemd/user/wretch-reconcile.service
# Runs one reconcile pass: transcribe any landed-but-unprocessed Craig zips, then
# cascade the downstream rebuild (wiki + podcast). Triggered by the .path unit
# (on folder change) and/or cron (recovery heartbeat). Safe to fire eagerly — the
# reconciler gates on `unzip -t` and holds a single-flight lock.
[Unit]
Description=Faerrin wretch — transcribe new Craig recordings, then publish (wiki + podcast)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
TimeoutStartSec=infinity
OOMScoreAdjust=-500
WorkingDirectory=/ruby/data/experiments/faerrin
EnvironmentFile=/ruby/data/experiments/faerrin/pkg/wretch/.env
ExecStart=/bin/bash -lc 'bun run --filter @faerrin/wretch process'
